# ansible/roles/supervisor/templates/message-consumer.ini
[program:{{ app_message_consumers_name }}]
directory = {{ app_root }}
command = php bin/console console messenger:consume async
process_name = %(program_name)s_%(process_num)02d
numprocs = {{ app_command_consumers_number }}
autostart = true
user = {{ system_user }}
